# filelist
/root/filelist
/etc/init.d/brook
/etc/init.d/frps
/etc/init.d/shadowsocks-go
/etc/nginx/1_www.srgb.xyz_bundle.crt
/etc/nginx/2_www.srgb.xyz.key
/etc/nginx/sites-enabled/default
/etc/nginx/sites-enabled/frps
/etc/nginx/sites-enabled/http
/etc/nginx/sites-enabled/https
/etc/shadowsocks-go/config.json
/root/.bashrc
/root/.vimrc
/root/bbr.sh
/root/fuck_net
/usr/bin/shadowsocks-server
/usr/local/brook/brook
/usr/local/brook/brook.conf
/usr/local/frps/frps
/usr/local/frps/frps.ini

